# jimmy-discord.service
[Unit]
Description=Jimmy Discord BitNet AI Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__WORKDIR__
ExecStart=__WORKDIR__/.venv/bin/python bot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jimmy-discord

[Install]
WantedBy=multi-user.target
